# service/apache/run
#!/bin/bash

sv start apache-confd || exit 1
exec apache2ctl -D FOREGROUND 
