# service/php-fpm/run
#!/bin/bash

exec php5-fpm -c /etc/php/ -y /etc/php/php-fpm.conf -F
